# deploy/reddzit-ssr.service.example
[Unit]
Description=Reddzit SSR meta tag server
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/reddzit-refresh
Environment=NODE_ENV=production
Environment=PUBLIC_BASE_URL=https://reddzit.seojeek.com
ExecStart=/usr/bin/node server/index.mjs
Restart=on-failure
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target

